# scripts/pi/lite/ariopi-signage.service
[Unit]
Description=ArioPi Player (MPV)
After=local-fs.target

[Service]
Type=simple
User=pi
Group=video
# GUI yok; MPV DRM/rpi ile doğrudan HDMI
Environment=DISPLAY=
WorkingDirectory=/home/pi
ExecStart=/usr/bin/python3 /opt/ariopi-signage/signage_client.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
